# xfce/udev-rules/10-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/bin/sh -c 'cd /home/hbotha/scripts; /bin/sh monitor-layout.sh'"
